# debian/salt-api.service
[Unit]
Description=The Salt Api
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/salt-api

[Install]
WantedBy=multi-user.target
